# deploy/forge-simple.service
[Unit]
Description=Forge Simple backend (prompt2app)
After=network.target

[Service]
EnvironmentFile=/etc/forge-simple.env
WorkingDirectory=/opt/forge-simple
ExecStart=/opt/forge-simple/.venv/bin/python -m backend.server
Restart=always
RestartSec=5

# Hardening DISABLED — каждая опция блокировала bundled Claude CLI на initialize.
# TODO: восстановить по одной с тестированием когда определим точный набор.

[Install]
WantedBy=multi-user.target
